# Ansible/project07/vprofile/ansible/roles/weblogic/templates/weblogic.service.j2
[Unit]
Description=Oracle WebLogic Server - {{ weblogic_domain_name }}
After=network.target
Documentation=https://docs.oracle.com/en/middleware/standalone/weblogic-server/
StartLimitIntervalSec=300
StartLimitBurst=3

[Service]
Type=forking
User={{ weblogic_user }}
Group={{ weblogic_group }}
WorkingDirectory={{ domain_home }}

# Service commands
ExecStart={{ domain_home }}/bin/daemon-start.sh
ExecStop={{ domain_home }}/bin/daemon-stop.sh

# PID management
PIDFile={{ domain_home }}/AdminServer.pid

# Timeouts
TimeoutStartSec=1500
TimeoutStopSec=300

# Restart policy
Restart=on-failure
RestartSec=60

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=weblogic-{{ weblogic_domain_name }}

[Install]
WantedBy=multi-user.target
